# supervisord-haproxy.conf
[program:haproxy]
command=/start-haproxy.sh
numprocs=1
autostart=true
autorestart=true